# .config/systemd/user/logitech-devices.target
[Unit]
Description=Logitech devices are ready
